# games-util/capsule-farmer-evolved/files/capsule-farmer-evolved.service
[Unit]
Description=CapsuleFarmerEvolved
After=network.target

[Service]
ExecStart=/usr/bin/capsule-farmer-evolved
Restart=always
RestartSec=60